# provisioning/hdmi-control/hdmi-control.service
[Unit]
Description=HDMI Control

[Service]
Type=oneshot
ExecStart=/usr/bin/hdmi-control

[Install]
WantedBy=multi-user.target
